# scripts/homepage.cacheUpdate.timer
[Unit]
Description=run homepage.cacheUpdate each 15 minutes

[Timer]
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target