# scripts/homepage.watcher
#!/bin/sh
### BEGIN INIT INFO
# Provides:          homepage.watcher
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      1 0 6
# Short-Description: ...
# Description: ...
### END INIT INFO

while true; do
        curl "http://olehkrupko.com/feedUpdate/cache";
        sleep 10m;
done;
