# scripts/homepage.cacheUpdate.timer
[Unit]
Description=run homepage.cacheUpdate each 15 minutes

After=network.target
Wants=network-online.target
After=network-online.target
After=homepage.server.launch.service
Requires=homepage.cacheUpdate.service

[Timer]
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target