# scripts/homepage.server.reload.timer
[Unit]
Description=homepage.server.reload.timer: run homepage.server.reload each 12 hours

After=network.target
Wants=network-online.target
After=network-online.target
After=homepage.server.launch.service

[Timer]
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target