# scripts/homepage.cacheUpdate.all.timer
[Unit]
Description=homepage.cacheUpdate.all.timer: run homepage.cacheUpdate.all each 84 minutes

[Timer]
OnUnitActiveSec=84min

[Install]
WantedBy=timers.target