# scripts/homepage.watcher
#!/bin/sh

### BEGIN INIT INFO
# Provides:          com.olehkrupko.homepage.watcher
# Required-Start:
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      1 0 6
# Short-Description: launches Django-cacher demon
### END INIT INFO

while true
do
        curl "http://home.olehkrupko.com/bookUpdates/cache";
        sleep 10m
done